# tools/load-balancer/etc/supervisor.d/haproxy.conf
[program:haproxy]
command=bin/haproxy-wrapper haproxy -p var/haproxy.pid -f etc/haproxy.conf -f var/haproxy-hue.conf
autostart=true
autorestart=true
stopsignal=USR1
