# tools/load-balancer/etc/supervisor.d/nginx.conf
[program:nginx]
command=nginx -p %(here)s/../ -c etc/nginx.conf
autostart=true
autorestart=false
